# deploy/systemd/nr-timetable-index-prebuild.timer
[Unit]
Description=Fallback daily run of NR timetable index prebuild

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10m
Unit=nr-timetable-index-prebuild.service

[Install]
WantedBy=timers.target
